# lib-messagebus/platforms/kafka/kafka.service
[Unit]
Requires=zookeeper.service
After=zookeeper.service

[Service]
Type=simple
User=kafka
ExecStart=/bin/sh -c '/home/kafka/kafka_install/bin/kafka-server-start.sh /home/kafka/kafka_install/config/server.properties > /home/kafka/kafka_install/kafka.log 2>&1'
ExecStop=/home/kafka/kafka_install/bin/kafka-server-stop.sh
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
